# deploy/opus_aiops.service
[Unit]
Description=Opus AI Ops Service
After=network.target opus_watchdog.service
Wants=network-online.target

[Service]
Type=simple
User=aiops
Group=aiops
WorkingDirectory=/var/www
EnvironmentFile=/var/www/opus_aiops/.env
ExecStart=/var/www/aiops-venv/bin/uvicorn opus_aiops.main:app --host 127.0.0.1 --port 9100 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/www/logs/aiops/aiops.log
StandardError=append:/var/www/logs/aiops/aiops.error.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/www/opus_aiops /var/www/aiops-venv
ReadWritePaths=/var/www/opus_aiops/aiops.db /var/www/opus_aiops/aiops.db-wal /var/www/opus_aiops/aiops.db-shm /var/www/logs/aiops
# Read-only access to trader storage for file-based evidence
ReadOnlyPaths=/var/www/storage

[Install]
WantedBy=multi-user.target
