# deploy/opus_aiops.service
[Unit]
Description=Opus AI Ops Service
After=network.target opus_watchdog.service
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/var/www
EnvironmentFile=/var/www/opus_aiops/.env
ExecStart=/var/www/aiops-venv/bin/uvicorn opus_aiops.main:app --host 127.0.0.1 --port 9100 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/www/logs/aiops/aiops.log
StandardError=append:/var/www/logs/aiops/aiops.error.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
